# dashcam-crawler.service
[Unit]
Description=Dashcam Crawler
After=network.target
StartLimitBurst=10
StartLimitIntervalSec=120
StartLimitAction=reboot

[Service]
EnvironmentFile=/etc/dashcam-crawler.conf
Type=exec
User=${USER}
WorkingDirectory=${DATA_DIR}
Environment="PYTHONPATH=${PROJECT_ROOT}"
ExecStart=${VENV_DIR}/bin/python3 -m crawler.main
Restart=always
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
