# dashcam-crawler-restart.service
[Unit]
Description=Restart Dashcam Crawler after config changes

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart dashcam-crawler.service
