# dashcam-crawler-restart.path
[Unit]
Description=Watch Dashcam Crawler config changes

[Path]
PathChanged=/etc/dashcam-crawler.conf
Unit=dashcam-crawler-restart.service

[Install]
WantedBy=multi-user.target
