# debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        if [ -d /run/systemd/system ]; then
            systemctl disable --now dashcam-crawler-restart.path 2>/dev/null || true
            systemctl disable --now dashcam-crawler.service 2>/dev/null || true
        fi
        ;;
    failed-upgrade)
        ;;
    *)
        echo "prerm called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
